# sessions/prometheus-omero-py-metrics.service
[Unit]
Description=prometheus-omerosessions

[Service]
User=omero-prometheus
EnvironmentFile=-/etc/sysconfig/prometheus-omerosessions
Environment="PYTHONPATH=/opt/omero/server/OMERO.server/lib/python"
ExecStart=/opt/omero-prometheus-tools/venv/bin/python /opt/omero-prometheus-tools/src/sessions/metrics.py $OPTIONS
RestartSec=10
Restart=on-failure
SyslogIdentifier=prometheus-omerosessions

[Install]
WantedBy=multi-user.target
